# copy/webapi.conf
[program:webapi]
directory=root/webapi
command=/root/webapi/WebApi
autorestart=true
autostart=true